# deploy/neli.service
[Unit]
Description=Neli Muhendislik Backend Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/nelimuhendislik
ExecStart=/root/.bun/bin/bun run server/index.ts
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3001
EnvironmentFile=/var/www/nelimuhendislik/.env

[Install]
WantedBy=multi-user.target
